# ops/minibox.slice
[Unit]
Description=Minibox slice

[Slice]
Delegate=yes
TasksMax=infinity
