# ops/miniboxd.service
[Unit]
Description=Minibox daemon
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'if [ -w /sys/fs/cgroup/minibox.slice/miniboxd.service/cgroup.subtree_control ]; then echo "+cpu +memory +pids" > /sys/fs/cgroup/minibox.slice/miniboxd.service/cgroup.subtree_control; fi'
ExecStart=/usr/local/bin/miniboxd
Restart=on-failure
RestartSec=2
User=root
Group=root
Slice=minibox.slice
Delegate=yes
Environment=MINIBOX_CGROUP_ROOT=/sys/fs/cgroup/minibox.slice/miniboxd.service/minibox
RuntimeDirectory=minibox
RuntimeDirectoryMode=0700
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
